# examples/gentoo-soaptube.in
#!/sbin/runscript

user=liquidsoap
prefix=@prefix@
exec_prefix=@exec_prefix@
confdir=@sysconfdir@/liquidsoap
soaptube=@bindir@/soaptube
logdir=@localstatedir@/log/liquidsoap
rundir=@localstatedir@/run/liquidsoap

depend() {
  after net
}

start() {
  einfo "Starting relays..."
  cd $confdir
  if [ -f *.out.conf ] ; then
    for out in *.out.conf ; do
      log=${out%.conf}.log
      ebegin " $out"
      sudo -u $user $soaptube -d -l ${confdir}/$out \
        > $logdir/${log}.1 2> $logdir/${log}.2 &
      eend $?
    done
  fi
}

stop() {
  ebegin "Stopping relay(s)"
  cd ${rundir}
  if [ -f *.out.pid ] ; then
    for liq in *.out.pid ; do
      kill `cat $liq`
      rm $liq
    done
  fi
  eend 0
}

restart() {
  svc_stop
  einfo Sleeping 3 seconds...
  sleep 3
  svc_start
}
